# deploy/crontab
# CarreiraHub - root user crontab
#
# Active app jobs are systemd timers:
#   systemctl list-timers 'carreirahub-cron@*'
#
# Keep the root crontab limited to the Telegram watcher. Do not duplicate the
# app cron jobs here while systemd timers are enabled.

SHELL=/bin/bash
MAILTO=""

*/10 * * * * /opt/carreirahub/cron-watch.sh >> /var/log/carreirahub-watch.log 2>&1
